# efixo/efixo-base-files/files/nb6-main/etc/init.d/ath
#!/bin/sh /etc/rc.common

START=15

boot() {
	chmod +s /usr/bin/ath*
}
