# efixo/efixo-base-files/files/main/etc/init.d/sipalg
#!/bin/sh /etc/rc.common

logger "[sipalg] $$ $action"

start() {
	insmod nf_nat_sip
}

stop() {
	rmmod nf_nat_sip
}

